# deploy/moderator-bot.service
[Unit]
Description=Moderator Bot (Tech Support Telegram Bot)
After=network.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/telegram-knowledge-bot
ExecStart=/opt/telegram-knowledge-bot/venv/bin/python run_moderator_bot.py
Restart=always
RestartSec=5
EnvironmentFile=/opt/telegram-knowledge-bot/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=moderator-bot

[Install]
WantedBy=multi-user.target
