# deploy/ustaz-bot.service
[Unit]
Description=Ustaz Bot (Ustaz Telegram Bot)
After=network.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/telegram-knowledge-bot
ExecStart=/opt/telegram-knowledge-bot/venv/bin/python run_ustaz_bot.py
Restart=always
RestartSec=5
EnvironmentFile=/opt/telegram-knowledge-bot/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ustaz-bot

[Install]
WantedBy=multi-user.target
